# service/weclaw.service
[Unit]
Description=WeClaw - WeChat AI Agent Bridge
After=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/local/bin/weclaw start -f
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
